# system/yunapbx-appliance/package_files/etc/hotplug.d/iface/10-openvpn
#!/bin/sh

if [ "${ACTION}" = "ifup" ]; then
    if [ "${INTERFACE}" = "lan" ] || [ "${INTERFACE}" = "lan6" ]; then
        /etc/init.d/openvpn restart
    fi
fi
