# net/yunapbx-appliance/image_files/etc/init.d/asterisk
#!/bin/sh /etc/rc.common
# Copyright (C) 2014 OpenWrt.org

START=99

start() {
       mkdir -p /var/log/asterisk
       mkdir -p /var/spool/asterisk

       chown -R asterisk:asterisk /var/log/asterisk
       chown -R asterisk:asterisk /var/spool/asterisk
       chown -R asterisk:asterisk /usr/lib/asterisk

       service_start /usr/sbin/asterisk
}

stop() {
       /usr/sbin/asterisk -rx 'core stop now'
}

reload() {
       /usr/sbin/asterisk -rx 'core reload'
}
