# net/yunapbx-appliance/package_files/etc/rc.d/S45firstboot
#!/bin/sh /etc/rc.common

start() {
    if [ ! -f "/etc/ssl/private/server.key" ] || [ ! -f "/etc/ssl/private/server.crt" ]; then
        /usr/sbin/sslkeygen
    fi

    rm /etc/rc.d/S45firstboot
}
