# net/yunapbx-appliance/package_files/etc/rc.d/S96firstboot
#!/bin/sh /etc/rc.common

start() {
    /usr/bin/php-cli /usr/share/yunapbx/db/install_upgrade_db.php
    
    rm /etc/rc.d/S96firstboot
    
}
