# system/yunalbl-appliance/package_files/etc/hotplug.d/iface/10-zabbix_agentd
#!/bin/sh

if [ "${ACTION}" = "ifup" ]; then
    if [ "${INTERFACE}" = "lan" ] || [ "${INTERFACE}" = "lan6" ]; then
        /etc/init.d/zabbix_agentd restart
    fi
fi
